# ops/com.signalforge.kill_switch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.signalforge.kill_switch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/varunteja/SignalForge/.venv/bin/python</string>
        <string>/Users/varunteja/SignalForge/scripts/kill_switch.py</string>
        <string>--log-level</string>
        <string>INFO</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/varunteja/SignalForge</string>

    <!-- Runs every hour at :30, 29 minutes AFTER the paper trader,
         so it reads a journal that includes the most recent fills. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/varunteja/SignalForge/fund_data/logs/kill_switch.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/varunteja/SignalForge/fund_data/logs/kill_switch.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
